# ha_sentinel/rootfs/etc/services.d/ha-sentinel/run
#!/usr/bin/with-contenv bashio
exec python3 /app/ha_sentinel.py
